# recipes-extended/canservice/files/can.service
[Unit]
Description=CAN startup service
After=network.target

[Service]
ExecStart=/usr/bin/start-can

[Install]
WantedBy=multi-user.target
